# misc/cron.d/dnssec
0 1 * * 6 root /usr/local/autodnssec/dnssec.sh --re-sign
*/20 * * * * root /usr/local/autodnssec/dnssec.sh --check-dnssec
